# scripts/tests/supervisor/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false  # 容器中必须前台运行

# 开启 TCP 监听（允许远程访问）
[inet_http_server]
port=0.0.0.0:9001  # 监听所有网卡的 9001 端口（0.0.0.0 表示允许外部访问）
username=admin     # 可选：远程访问的用户名（需配合密码）
password=123456    # 可选：远程访问的密码

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# 配置 supervisorctl 远程连接（客户端用）
[supervisorctl]
serverurl=http://0.0.0.0:9001  # 与上面的 TCP 端口一致
username=admin                 # 与上面的用户名一致
password=123456                # 与上面的密码一致

# 被管理的进程（示例）
[program:mybin1]
command=/usr/bin/python3 /app/v1/mybin1
autostart=true
autorestart=true

[program:mybin2]
command=/usr/bin/python3 /app/v1/mybin2
autostart=true
autorestart=true

[program:mybin3]
command=/usr/bin/python3 /app/v1/mybin3
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"  ; nginx 需前台运行（关闭自身守护进程模式）
autostart=true
autorestart=true